# omabutton.conf
[program:omabutton]
command=python /usr/local/omabutton/player.py
autostart=true
autorestart=true
stderr_logfile=/var/log/omabutton.log
stdout_logfile=/var/log/omabutton.err.log